# services/fw-fanctrl-suspend.service
[Unit]
Description=Framework Fan Controller sleep hook
Before=sleep.target
StopWhenUnneeded=yes
Requires=fw-fanctrl.service
After=fw-fanctrl.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart="%PYTHON_SCRIPT_INSTALLATION_PATH%" pause
ExecStop="%PYTHON_SCRIPT_INSTALLATION_PATH%" resume

[Install]
WantedBy=sleep.target
